# systemd/github-agent.service
[Unit]
Description=GitHub Actions Webhook Deploy Agent
After=network.target
Wants=network.target

[Service]
Type=simple
User=github-agent
Group=github-agent

# Secret loaded from a file with mode 0600 owned by github-agent
EnvironmentFile=/etc/github-agent/env

ExecStart=/usr/local/bin/github-agent -config /etc/github-agent/config.yaml
Restart=on-failure
RestartSec=5s

# --- Hardening (deliberately light) ---
# Deploy scripts run via sudo INSIDE this service's mount namespace, so:
#   - ProtectHome / ProtectSystem=strict would hide /home from the scripts
#   - SystemCallFilter / ProtectKernel* / CapabilityBoundingSet= imply
#     NoNewPrivileges=yes on non-root services, which breaks sudo entirely
ProtectSystem=full
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
