# systemd/github-agent.service
# REFERENCE TEMPLATE ONLY. The real unit is generated by install.sh (write_systemd_unit)
# from config.yaml, with one ReadWritePaths= line per deploy_path and one
# BindReadOnlyPaths= per deploy_user's .docker. If you install this file by hand you MUST
# add those lines yourself — replace the <DEPLOY_PATH> / <DEPLOY_USER_HOME> placeholders.

[Unit]
Description=GitHub Actions Webhook Deploy Agent
After=network.target
Wants=network.target

[Service]
Type=simple
User=github-agent
Group=github-agent

# Secret loaded from a file with mode 0600 owned by github-agent
EnvironmentFile=/etc/github-agent/env

ExecStart=/usr/local/bin/github-agent -config /etc/github-agent/config.yaml
Restart=on-failure
RestartSec=5s

# --- Hardening ---
# ProtectSystem=strict makes the whole FS read-only; ProtectHome=true replaces /home,
# /root and /run/user with empty tmpfs. The agent chdir's into deploy_path as the
# github-agent user BEFORE sudo escalates, so each deploy_path must be re-exposed
# read-write or that chdir fails with "chdir ...: permission denied".
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/github-agent
# One line per project — replace with your real deploy_path(s):
ReadWritePaths=<DEPLOY_PATH>

# docker reads registry credentials from the deploy_user's $HOME/.docker/config.json.
# ProtectHome=true hides it, so private-registry pulls fail with "pull access denied".
# Re-expose each deploy_user's .docker read-only. Leading '-' = skip if absent.
BindReadOnlyPaths=-<DEPLOY_USER_HOME>/.docker

# sudo relies on the setuid bit to escalate to the deploy_user. NoNewPrivileges=true
# blocks that (setuid is denied), breaking every deploy — it MUST stay false.
NoNewPrivileges=false

# Restrict syscalls to a safe set
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Isolate private /tmp
PrivateTmp=true

# Prevent loading kernel modules / tuning the kernel
ProtectKernelModules=true
ProtectKernelTunables=true

StandardOutput=append:/var/log/github-agent/agent.log
StandardError=append:/var/log/github-agent/agent.log

[Install]
WantedBy=multi-user.target
